# backend/deploy/shoresafe-backend.service
[Unit]
Description=ShoreSafe Backend API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/shoresafe/backend
EnvironmentFile=/etc/shoresafe-backend.env
ExecStart=/opt/shoresafe/backend/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=on-failure
User=shoresafe
Group=shoresafe

[Install]
WantedBy=multi-user.target
